# crontab
5 * * * * /usr/bin/bash /home/data_mine/rtmarchive/rtmarchive.sh >> /tmp/rtmarchive.log
CRON_TZ=UTC
0 1 * * * /usr/bin/bash /home/data_mine/rtmarchive/archivemgmt.sh >> /tmp/archivemgmt.log

