# rpmbuild/BUILD/rtmarchive-1.2.0/opt/rtmarchive/cron/rtmarchive.crontab
0 0 * * * /usr/bin/bash /opt/rtmarchive/queryrumc.sh -u 2>&1 | sed -e "s/^/$(date -Is) /" >> /var/log/rtmarchive/queryrumc.log
5 * * * * /usr/bin/bash /opt/rtmarchive/rtmarchive.sh 2>&1 | sed -e "s/^/$(date -Is) /" >> /var/log/rtmarchive/rtmarchive.log
CRON_TZ=UTC
0 1 * * * /usr/bin/bash /opt/rtmarchive/archivemgmt.sh 2>&1 | sed -e "s/^/$(date -Is) /" >> /var/log/rtmarchive/archivemgmt.log
0 3 * * * /usr/bin/bash /opt/rtmarchive/archivemgmtindex.sh 2>&1 | sed -e "s/^/$(date -Is) /" >> /var/log/rtmarchive/archiveindex.log

